# tech_bloggers/deployment/gunicorn.service
# Systemd service file for Gunicorn
# File location: /etc/systemd/system/gunicorn.service
#
# After creating this file:
# 1. Copy to /etc/systemd/system/gunicorn.service
# 2. Reload systemd: sudo systemctl daemon-reload
# 3. Enable service: sudo systemctl enable gunicorn
# 4. Start service: sudo systemctl start gunicorn
# 5. Check status: sudo systemctl status gunicorn

[Unit]
Description=Gunicorn daemon for Tech-In-Bytes Django application
After=network.target

[Service]
Type=notify
# User and group to run as (REPLACE 'django' with your user if different)
User=django
Group=www-data

# Working directory
WorkingDirectory=/var/www/techinbytes

# Environment variables
Environment="DJANGO_ENV=production"
Environment="PATH=/var/www/techinbytes/venv/bin"
Environment="PYTHONPATH=/var/www/techinbytes"

# Load additional environment variables from file
EnvironmentFile=/var/www/techinbytes/.env

# Gunicorn execution
ExecStart=/var/www/techinbytes/venv/bin/gunicorn \
    --config /var/www/techinbytes/gunicorn_config.py \
    tech_bloggers.wsgi:application

# Restart policy
Restart=on-failure
RestartSec=5s

# Security
# Prevent service from accessing home directories
ProtectHome=true
# Create private /tmp
PrivateTmp=true
# Prevent writing to /usr, /boot, /efi
ProtectSystem=strict
# Allow write access to these directories
ReadWritePaths=/var/www/techinbytes/media
ReadWritePaths=/var/www/techinbytes/logs
ReadWritePaths=/var/www/techinbytes/staticfiles

# Process limits
LimitNOFILE=65536
LimitNPROC=4096

# Timeout for service startup
TimeoutStartSec=60

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target

