# tech_bloggers/deployment/gunicorn.service
[Unit]
Description=Gunicorn daemon for Tech-In-Bytes Django application
After=network.target

[Service]
Type=notify
User=django
Group=django

# Working directory
WorkingDirectory=/var/www/techinbytes/tech_bloggers

# Runtime directory (creates /run/gunicorn/ automatically)
RuntimeDirectory=gunicorn

# Environment variables
Environment="DJANGO_ENV=production"
Environment="PATH=/var/www/techinbytes/venv/bin"
Environment="PYTHONPATH=/var/www/techinbytes/tech_bloggers"

# Load .env file
EnvironmentFile=/var/www/techinbytes/.env

# Gunicorn execution (socket in proper location)
ExecStart=/var/www/techinbytes/venv/bin/gunicorn \
    --bind unix:/run/gunicorn/gunicorn.sock \
    --workers 3 \
    --timeout 30 \
    --access-logfile /var/www/techinbytes/logs/gunicorn-access.log \
    --error-logfile /var/www/techinbytes/logs/gunicorn-error.log \
    --log-level info \
    tech_bloggers.wsgi:application

# Restart policy
Restart=on-failure
RestartSec=5s

# Security
ProtectHome=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/techinbytes/media
ReadWritePaths=/var/www/techinbytes/logs
ReadWritePaths=/var/www/techinbytes/staticfiles
ReadWritePaths=/home/django/.postgresql

# Process limits
LimitNOFILE=65536
LimitNPROC=4096

# Timeout
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target